# config/systemd/kitchensync.container
[Unit]
Description=Kitchen Sync

[Container]
Image=ghcr.io/thristhart/kitchensync:latest
Mount=type=bind,source=/srv/containers/kitchensync/db/,target=/mnt/db
Environment=DATABASE_URL=sqlite:/mnt/db/prod.sqlite3
PublishPort=localhost:4000:3000

[Install]
WantedBy=default.target